# etc/xinetd.d/rsyncd
service rsync
{
	socket_type = stream
	protocol    = tcp
	wait        = no
	user        = root
	server      = /home/core/gentoo/usr/bin/rsync
	server_args = --daemon
	disable     = yes
}
